# systemd/openkoutsi-frontend@.service
# LEGACY: bare-metal deploy unit, kept for reference only.
# The supported deploy model is now the GHCR container image
# (ghcr.io/openkoutsi/openkoutsi-web) pulled and run on the VM. This unit ran the
# rsynced `.next/standalone` build produced by the removed scripts/deploy-frontend.sh.

[Unit]
Description=openkoutsi frontend
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/projects/openkoutsi-web/.next/standalone
Environment=PORT=3000
# Adjust the node path if you switch Node versions via nvm.
ExecStart=/home/%i/.nvm/versions/node/v25.2.1/bin/node server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
