# systemd/openkoutsi-frontend@.service
[Unit]
Description=openkoutsi frontend
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/projects/openkoutsi-web/.next/standalone
Environment=PORT=3000
# Adjust the node path if you switch Node versions via nvm.
ExecStart=/home/%i/.nvm/versions/node/v25.2.1/bin/node server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
